# files/etc/supervisor.d/bird.ini
[program:bird]
command=/usr/sbin/bird -f -c/etc/bird.d/bird
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
